# root/etc/services.d/jellyfin/run
#!/usr/local/bin/execlineb -P
/usr/local/bin/s6-setuidgid bsd
/usr/local/bin/dotnet /usr/local/jellyfin/Jellyfin.Server.dll
    --datadir /config
    --cachedir /cache
    --ffmpeg /usr/local/bin/ffmpeg
    --webdir /usr/local/jellyfin/jellyfin-web